# rootfs/etc/init.d/game-ram-setup
#!/sbin/openrc-run
# Populate Shell RPG runtime under tmpfs (/tmp) at boot.

depend() {
	need localmount
	after localmount
}

start() {
	ebegin "Shell RPG RAM setup"

	rm -rf /tmp/game_map /tmp/bin /tmp/player.json /tmp/lib /tmp/usr

	if [ ! -x /usr/local/share/shell_rpg_engine_mpy/install.sh ]; then
		ewend 1 "install.sh missing"
		return 1
	fi

	if ! su user42 -s /bin/bash -c \
		'cd /usr/local/share/shell_rpg_engine_mpy && ./install.sh engine.py' \
		>>/tmp/game-ram-setup.log 2>&1; then
		ewend 1 "install.sh failed (see /tmp/game-ram-setup.log)"
		return 1
	fi

	if ! [ -d /tmp/game_map/village ] || ! [ -d /tmp/game_map/montagne ] || \
		! [ -r /tmp/player.json ] || ! [ -x /tmp/bin/ls ]; then
		ewend 1 "game map incomplete"
		return 1
	fi

	eend 0
}

stop() {
	need_stop=no
}
